# basket-manager.service
[Unit]
Description=Solana Basket Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=marc
WorkingDirectory=/home/marc/arb-agent
EnvironmentFile=/home/marc/arb-agent/.env
ExecStart=/usr/bin/node /home/marc/arb-agent/dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=basket-manager

# Harden
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/marc/arb-agent

[Install]
WantedBy=multi-user.target
